# dev-ops/systemd/vitrobud-opencads-backup.service.in
[Unit]
Description=Back up Open CAD Studio private runtime configuration

[Service]
Type=oneshot
WorkingDirectory=@PROJECT_DIR@
ExecStart=@PROJECT_DIR@/dev-ops/backup
ExecStart=@PROJECT_DIR@/dev-ops/cleanup
